# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:ollama]
command=ollama serve
priority=1
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/ollama.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3

[program:model-pull]
command=/bin/bash /app/pull-model.sh
priority=5
autorestart=false
redirect_stderr=true
stdout_logfile=/var/log/model-pull.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=1

[program:fastapi]
command=/app/venv/bin/uvicorn sd-server.app:app --host 0.0.0.0 --port 80 --workers 1
directory=/app
priority=10
environment=PATH="/app/venv/bin:%(ENV_PATH)s"
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/fastapi.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=3
startsecs=10